# laravel-worker.conf
[program:ImportSales]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work  --queue=ImportSales
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/laravel-queue.log
stopwaitsecs=3600